# pihub/supervisord.conf


{% if settings.FCGI %}

[program:runserver]
exclude=true

[program:autoreload]
exclude=true

[program:fcgi]
command={{ PYTHON }} {{ PROJECT_DIR }}/manage.py runfcgi host=127.0.0.1 port={{ settings.FCGI_PORT }} daemonize=false

{% else %}

[program:runserver]
autoreload=false
command={{ PYTHON }} {{ PROJECT_DIR }}/manage.py runserver 8020

{% endif %}



[program:celerycam]
command={{ PYTHON }} {{ PROJECT_DIR }}/manage.py celerycam 


[program:celeryd]
command={{ PYTHON }} {{ PROJECT_DIR }}/manage.py celeryd --settings={{ environ.BDB_SETTINGS }} -E --loglevel={{ settings.CELERY_LOG_LEVEL }}

directory={{ PROJECT_DIR }}
numprocs=1
startsecs=5

; Need to wait for currently executing tasks to finish at shutdown.
; Increase this if you have very long running tasks.
stopwaitsecs = 600

; if rabbitmq is supervised, set its priority higher
; so it starts first
priority=998
